# roles/undervolt/files/opt/undervolt/undervolt.service
[Unit]
Description=undervolt

[Service]
Type=oneshot
ExecStart=/opt/undervolt/undervolt.py -v --core -55 --gpu -55 --cache -55 --uncore -55
